# docker/etc/s6-overlay/s6-rc.d/caddy/run
#!/command/execlineb -P
s6-setuidgid nextjs
cd /app
caddy run Caddyfile
