# docker/etc/s6-overlay/s6-rc.d/caddy/run
#!/command/execlineb -P
cd /app
caddy run Caddyfile
