# docker/etc/s6-overlay/s6-rc.d/next/run
#!/command/execlineb -P
cd /app
export PORT 3500
s6-setuidgid nextjs
node server.js
